# kubernetes/dev/sphinx/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f

[program:sphinx]
command=/usr/bin/searchd -c /etc/sphinxsearch/sphinx.conf --nodetach
